# rc.d/rc.local
#!/bin/sh
#
# $Id: rc.local,v 1.4 2000/11/16 16:01:17 baggins Exp $
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.


# This must be last line !
# vi:syntax=sh:tw=78:ts=8:sw=4
