# rc.d/init.d/allowlogin
#!/bin/sh
#
# allowlogin	Allow users to login
#
# chkconfig: 2345 99 1
# description: Allow users to login.
#

# Source function library.
. /etc/rc.d/init.d/functions

[ ! -f /etc/sysconfig/system ] && exit 0

. /etc/sysconfig/system

# See how we were called.
case "$1" in
  start)
	if [ -f /etc/nologin.boot ]; then
		run_cmd "Allowing users to login" rm -f /etc/nologin /etc/nologin.boot_
	fi
        ;;
  stop)
	# /etc/nologin when shutting down system
	if [ "$DELAY_LOGIN" == "yes" ] && [ ! -f /etc/nologin ]; then
	        echo > /etc/nologin
	        echo "[EN] System shutdown in progress " >> /etc/nologin
	        echo "[PL] Zatrzymywanie systemu w toku" >> /etc/nologin
	        echo >> /etc/nologin
	        chmod 644 /etc/nologin
	        cp -fp /etc/nologin /etc/nologin.boot
	fi
        ;;
  *)
        echo "Usage: $0 {start|stop}"
        exit 1
esac

exit 0
