# rc.d/init.d/network
#!/bin/sh
#
# network	Bring up/down networking
#
# chkconfig:    2345 10 90
# description:  Activates/Deactivates all network interfaces configured to \
#		start at boot time.
#
# probe:	true

# $Id: network 5884 2005-04-25 16:01:57Z adamg $

PATH=/sbin:/usr/sbin:/bin:/usr/bin

# NLS
NLS_DOMAIN="rc-scripts"

# Source function library.
. /etc/rc.d/init.d/functions

if [ ! -f /etc/sysconfig/network ]; then
	nls "%s is missing. Can't continue." "/etc/sysconfig/network"
	exit 1
fi

. /etc/sysconfig/network

# Check that networking is up.
is_no "${NETWORKING}" && exit 0

# See how we were called.
case "$1" in
  start)
        touch /var/lock/subsys/network
        ;;
  stop)
        rm -f /var/lock/subsys/network
        ;;
  restart)
	$0 stop
	$0 start
	;;
  *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac

exit 0

# This must be last line !
# vi:syntax=sh:ts=8:sw=4
