# trunk/init/cryptsetup.conf
description "Disk encryption via cryptsetup"

start on pld.sysinit-done
stop on pld.shutdown-started

task

script
	. /etc/rc.d/init.d/cryptsetup
	init_crypto
end script
